# whizvoice.service
[Unit]
Description=WhizVoice FastAPI Application
After=network.target

[Service]
User=nginx
Group=nginx
WorkingDirectory=/var/www/whizvoice
ExecStart=/var/www/whizvoice/venv/bin/python -m uvicorn app:app --host 0.0.0.0 --port 8000 --workers 4
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=whizvoice

[Install]
WantedBy=multi-user.target 